# dnx_system/utils/services/dnx-syslog.service
[Unit]
Description=DNXFIREWALL Syslog Client Service
After=dnx-database-sqlite.service
StartLimitIntervalSec=3

[Service]
Type=simple
Restart=always
RestartSec=3
User=dnx
ExecStart=/usr/bin/python3 /home/dnx/dnxfirewall/dnx_run.py syslog
Environment=HOME_DIR=/home/dnx/dnxfirewall

[Install]
WantedBy=multi-user.target