# dnx_system/utils/services/dnx-ips.service
[Unit]
Description=DNXFIREWALL IPS Service
After=dnx-log.service
StartLimitIntervalSec=10

[Service]
Type=simple
Restart=always
RestartSec=3
StartLimitBurst=3
User=root
ExecStart=/usr/bin/python3 /home/dnx/dnxfirewall/dnx_run.py ips_ids
Environment=HOME_DIR=/home/dnx/dnxfirewall

[Install]
WantedBy=multi-user.target
