# services/dnx-startup.service
[Unit]
Description=DNXFIREWALL Startup Tasks

[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/python3 /home/dnx/dnxfirewall/dnx_configure/dnx_startup.py
Environment=HOME_DIR=/home/dnx/dnxfirewall

[Install]
WantedBy=network-pre.target
